# supervisord.conf
[supervisord]
nodeamon=true

[program:runserver]
command=python manage.py runserver 0.0.0.0:8000

[program:runbot]
command=python manage.py runbot
